# paperless-ngx/rootfs/etc/s6-overlay/s6-rc.d/init-addon/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Benoit Anastay Add-on: Paperless-ngx
# Run Parless-ngx
# ==============================================================================

cd /usr/src/paperless/src/ || exit

# Replace env
PAPERLESS_FILENAME_FORMAT=$(bashio::config 'filename')
PAPERLESS_OCR_LANGUAGE=$(bashio::config 'language')
PAPERLESS_OCR_LANGUAGES=$(bashio::config 'language_packages')
PAPERLESS_TIME_ZONE=$(bashio::info.timezone)
PAPERLESS_CONSUMER_ENABLE_BARCODES=$(bashio::config 'barcodes_enabled')
PAPERLESS_CONSUMER_ENABLE_ASN_BARCODE=$(bashio::config 'barcodes_asn')
PAPERLESS_CONSUMER_RECURSIVE=$(bashio::config 'consumer_recursive')
PAPERLESS_CONSUMER_SUBDIRS_AS_TAGS=$(bashio::config 'consumer_subdirs_as_tags')
PAPERLESS_CONSUMPTION_DIR=/share/paperless/consume
PAPERLESS_DATA_DIR=/config/data
PAPERLESS_MEDIA_ROOT=/share/paperless/media
PAPERLESS_USE_X_FORWARD_HOST=true
PAPERLESS_USE_X_FORWARD_PORT=true
PAPERLESS_CONSUMER_BARCODE_SCANNER=ZXING
PAPERLESS_CONFIGURATION_PATH="/config/paperless.conf"
PAPERLESS_HTTP_REMOTE_USER_HEADER_NAME="HTTP_X_REMOTE_USER_NAME"
PAPERLESS_TIKA_ENDPOINT=http://ca5234a0-tika-gotenberg:9998
PAPERLESS_TIKA_GOTENBERG_ENDPOINT=http://ca5234a0-tika-gotenberg:3000
PAPERLESS_TIKA_ENABLED=$(bashio::config 'tika_gotenberg')

if bashio::config.has_value 'usermap_uid'; then
  USERMAP_UID=$(bashio::config 'usermap_uid')
  echo -n "${USERMAP_UID}" > /var/run/s6/container_environment/USERMAP_UID
  export USERMAP_UID
fi

if bashio::config.has_value 'usermap_gid'; then
  USERMAP_GID=$(bashio::config 'usermap_gid')
  echo -n "${USERMAP_GID}" > /var/run/s6/container_environment/USERMAP_GID
  export USERMAP_GID
fi

#This set the URL of paperless
if bashio::config.has_value 'url'; then
  PAPERLESS_URL=$(bashio::config 'url')
  bashio::log.info "Use custom defined URL $PAPERLESS_URL"
else
  PAPERLESS_URL="http://"+$(bashio::info.hostname)+":8000"
fi

if bashio::config.has_value 'polling_interval'; then
  PAPERLESS_CONSUMER_POLLING=$(bashio::config 'polling_interval')
  bashio::log.info "Pooling interval set to $PAPERLESS_CONSUMER_POLLING"
else
  PAPERLESS_CONSUMER_POLLING=0
fi

export PAPERLESS_URL
export PAPERLESS_FILENAME_FORMAT
export PAPERLESS_OCR_LANGUAGE
export PAPERLESS_OCR_LANGUAGES
export PAPERLESS_TIME_ZONE
export PAPERLESS_CONSUMER_POLLING
export PAPERLESS_CONSUMPTION_DIR
export PAPERLESS_DATA_DIR
export PAPERLESS_MEDIA_ROOT
export PAPERLESS_CONSUMER_ENABLE_BARCODES
export PAPERLESS_CONSUMER_ENABLE_ASN_BARCODE
export PAPERLESS_CONSUMER_RECURSIVE
export PAPERLESS_CONSUMER_SUBDIRS_AS_TAGS
export PAPERLESS_USE_X_FORWARD_HOST
export PAPERLESS_USE_X_FORWARD_PORT
export PAPERLESS_CONSUMER_BARCODE_SCANNER
export PAPERLESS_CONFIGURATION_PATH
export PAPERLESS_HTTP_REMOTE_USER_HEADER_NAME
export PAPERLESS_TIKA_ENDPOINT
export PAPERLESS_TIKA_GOTENBERG_ENDPOINT
export PAPERLESS_TIKA_ENABLED

echo -n "${PAPERLESS_URL}" > /var/run/s6/container_environment/PAPERLESS_URL
echo -n "${PAPERLESS_FILENAME_FORMAT}" > /var/run/s6/container_environment/PAPERLESS_FILENAME_FORMAT
echo -n "${PAPERLESS_OCR_LANGUAGE}" > /var/run/s6/container_environment/PAPERLESS_OCR_LANGUAGE
echo -n "${PAPERLESS_OCR_LANGUAGES}" > /var/run/s6/container_environment/PAPERLESS_OCR_LANGUAGES
echo -n "${PAPERLESS_TIME_ZONE}" > /var/run/s6/container_environment/PAPERLESS_TIME_ZONE
echo -n "${PAPERLESS_CONSUMER_POLLING}" > /var/run/s6/container_environment/PAPERLESS_CONSUMER_POLLING
echo -n "${PAPERLESS_CONSUMPTION_DIR}" > /var/run/s6/container_environment/PAPERLESS_CONSUMPTION_DIR
echo -n "${PAPERLESS_DATA_DIR}" > /var/run/s6/container_environment/PAPERLESS_DATA_DIR
echo -n "${PAPERLESS_MEDIA_ROOT}" > /var/run/s6/container_environment/PAPERLESS_MEDIA_ROOT
echo -n "${PAPERLESS_CONSUMER_ENABLE_BARCODES}" > /var/run/s6/container_environment/PAPERLESS_CONSUMER_ENABLE_BARCODES
echo -n "${PAPERLESS_CONSUMER_ENABLE_ASN_BARCODE}" > /var/run/s6/container_environment/PAPERLESS_CONSUMER_ENABLE_ASN_BARCODE
echo -n "${PAPERLESS_CONSUMER_RECURSIVE}" > /var/run/s6/container_environment/PAPERLESS_CONSUMER_RECURSIVE
echo -n "${PAPERLESS_CONSUMER_SUBDIRS_AS_TAGS}" > /var/run/s6/container_environment/PAPERLESS_CONSUMER_SUBDIRS_AS_TAGS
echo -n "${PAPERLESS_USE_X_FORWARD_HOST}" > /var/run/s6/container_environment/PAPERLESS_USE_X_FORWARD_HOST
echo -n "${PAPERLESS_USE_X_FORWARD_PORT}" > /var/run/s6/container_environment/PAPERLESS_USE_X_FORWARD_PORT
echo -n "${PAPERLESS_CONSUMER_BARCODE_SCANNER}" > /var/run/s6/container_environment/PAPERLESS_CONSUMER_BARCODE_SCANNER
echo -n "${PAPERLESS_CONFIGURATION_PATH}" > /var/run/s6/container_environment/PAPERLESS_CONFIGURATION_PATH
echo -n "${PAPERLESS_HTTP_REMOTE_USER_HEADER_NAME}" > /var/run/s6/container_environment/PAPERLESS_HTTP_REMOTE_USER_HEADER_NAME
echo -n "${PAPERLESS_TIKA_ENDPOINT}" > /var/run/s6/container_environment/PAPERLESS_TIKA_ENDPOINT
echo -n "${PAPERLESS_TIKA_GOTENBERG_ENDPOINT}" > /var/run/s6/container_environment/PAPERLESS_TIKA_GOTENBERG_ENDPOINT
echo -n "${PAPERLESS_TIKA_ENABLED}" > /var/run/s6/container_environment/PAPERLESS_TIKA_ENABLED

# Load user-provided Paperless configuration if it exists
CONF_FILE="/config/paperless.conf"
if [ -f "$CONF_FILE" ]; then
    bashio::log.info "Loading environment overrides from $CONF_FILE"
    while IFS='=' read -r key value; do
        # skip comments and empty lines
        [ -z "$key" ] && continue
        case "$key" in \#*) continue ;; esac
        echo -n "$value" > "/var/run/s6/container_environment/$key"
    done < "$CONF_FILE"
fi
