# paperless-ngx/rootfs/etc/s6-overlay/s6-rc.d/paperless/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Benoit Anastay Add-on: Paperless-ngx
# Run Parless-ngx
# ==============================================================================

cd /usr/src/paperless/src/ || exit

# Replace env
PAPERLESS_FILENAME_FORMAT=$(bashio::config 'filename')
PAPERLESS_OCR_LANGUAGE=$(bashio::config 'language')
PAPERLESS_OCR_LANGUAGES=$(bashio::config 'language_packages')
PAPERLESS_ADMIN_USER=$(bashio::config 'default_superuser.username')
PAPERLESS_ADMIN_MAIL=$(bashio::config 'default_superuser.email')
PAPERLESS_ADMIN_PASSWORD=$(bashio::config 'default_superuser.password')
PAPERLESS_TIME_ZONE=$(bashio::info.timezone)
PAPERLESS_CONSUMER_ENABLE_BARCODES=$(bashio::config 'barcodes_enabled')
PAPERLESS_CONSUMER_ENABLE_ASN_BARCODE=$(bashio::config 'barcodes_asn')
PAPERLESS_CONSUMER_ASN_BARCODE_PREFIX==$(bashio::config 'barcodes_asn_prefix')
PAPERLESS_CONSUMER_BARCODE_UPSCALE==$(bashio::config 'barcodes_upscale')
PAPERLESS_CONSUMER_BARCODE_DPI==$(bashio::config 'barcodes_dpi')
PAPERLESS_CONSUMER_BARCODE_MAX_PAGES==$(bashio::config 'barcodes_max_pages')
PAPERLESS_CONSUMER_RECURSIVE=$(bashio::config 'consumer_recursive')
PAPERLESS_CONSUMER_SUBDIRS_AS_TAGS=$(bashio::config 'consumer_subdirs_as_tags')
PAPERLESS_CONSUMPTION_DIR=/share/paperless/consume
PAPERLESS_DATA_DIR=/config/data
PAPERLESS_MEDIA_ROOT=/share/paperless/media
PAPERLESS_USE_X_FORWARD_HOST=true
PAPERLESS_USE_X_FORWARD_PORT=true
PAPERLESS_CONSUMER_BARCODE_SCANNER=ZXING
PAPERLESS_CONFIGURATION_PATH="/config/paperless.conf"
PAPERLESS_HTTP_REMOTE_USER_HEADER_NAME="HTTP_X_REMOTE_USER_NAME"

if bashio::config.has_value 'usermap_uid'; then
  USERMAP_UID=$(bashio::config 'usermap_uid')
  echo -n "${USERMAP_UID}" > /var/run/s6/container_environment/USERMAP_UID
  export USERMAP_UID
fi

if bashio::config.has_value 'usermap_gid'; then
  USERMAP_GID=$(bashio::config 'usermap_gid')
  echo -n "${USERMAP_GID}" > /var/run/s6/container_environment/USERMAP_GID
  export USERMAP_GID
fi

#This set the URL of paperless
if bashio::config.has_value 'url'; then
  PAPERLESS_URL=$(bashio::config 'url')
  bashio::log.info "Use custom defined URL $PAPERLESS_URL"
else
  PAPERLESS_URL="http://"+$(bashio::info.hostname)+":8000"
fi

if bashio::config.has_value 'polling_interval'; then
  PAPERLESS_CONSUMER_POLLING=$(bashio::config 'polling_interval')
  bashio::log.info "Pooling interval set to $PAPERLESS_CONSUMER_POLLING"
else
  PAPERLESS_CONSUMER_POLLING=0
fi

export PAPERLESS_URL
export PAPERLESS_FILENAME_FORMAT
export PAPERLESS_OCR_LANGUAGE
export PAPERLESS_OCR_LANGUAGES
export PAPERLESS_ADMIN_USER
export PAPERLESS_ADMIN_MAIL
export PAPERLESS_ADMIN_PASSWORD
export PAPERLESS_TIME_ZONE
export PAPERLESS_CONSUMER_POLLING
export PAPERLESS_CONSUMPTION_DIR
export PAPERLESS_DATA_DIR
export PAPERLESS_MEDIA_ROOT
export PAPERLESS_CONSUMER_ENABLE_BARCODES
export PAPERLESS_CONSUMER_ENABLE_ASN_BARCODE
export PAPERLESS_CONSUMER_ASN_BARCODE_PREFIX
export PAPERLESS_CONSUMER_BARCODE_UPSCALE
export PAPERLESS_CONSUMER_BARCODE_DPI
export PAPERLESS_CONSUMER_BARCODE_MAX_PAGES
export PAPERLESS_CONSUMER_RECURSIVE
export PAPERLESS_CONSUMER_SUBDIRS_AS_TAGS
export PAPERLESS_USE_X_FORWARD_HOST
export PAPERLESS_USE_X_FORWARD_PORT
export PAPERLESS_CONSUMER_BARCODE_SCANNER
export PAPERLESS_CONFIGURATION_PATH
export PAPERLESS_HTTP_REMOTE_USER_HEADER_NAME

if bashio::fs.file_exists '/config/paperless.conf'; then
    dotenv -f /config/paperless.conf -e true run /sbin/docker-entrypoint.sh "/usr/local/bin/paperless_cmd.sh"
else
    /sbin/docker-entrypoint.sh "/usr/local/bin/paperless_cmd.sh"
fi
