# paperless-ngx/rootfs/etc/s6-overlay/s6-rc.d/paperless/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Benoit Anastay Add-on: Paperless-ngx
# Run Parless-ngx
# ==============================================================================

/sbin/docker-entrypoint.sh