# paperless-ngx/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Benoit Anastay Add-on: Paperlexx-ngx
# Configures NGINX for use with Paperless-ngx
# ==============================================================================

# Generate direct access configuration, if enabled.
if bashio::var.has_value "$(bashio::addon.port 80)"; then
    bashio::log.blue "Direct Access is enabled"
    bashio::config.require.ssl
    bashio::var.json \
        certfile "$(bashio::config 'certfile')" \
        keyfile "$(bashio::config 'keyfile')" \
        ssl "^$(bashio::config 'ssl')" \
        | tempio \
            -template /etc/nginx/templates/direct.gtpl \
            -out /etc/nginx/servers/direct.conf

fi

bashio::log.blue "Ingress is enabled"
# Generate Ingress NGINX configuration
bashio::var.json \
    interface "$(bashio::addon.ip_address)" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf