# paperless-ngx/rootfs/etc/services.d/supercronic/run
#!/command/with-contenv bashio
# shellcheck shell=bash
exec /usr/local/bin/supercronic /etc/cron.d/paperless-backup