# paperless-ngx/rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Benoit Anastay Add-on: Paperless-ngx
# Run redis
# ==============================================================================
sysctl vm.overcommit_memory=1

redis-server --daemonize yes