# paperless-ngx/rootfs/etc/s6-overlay/s6-rc.d/init-paperless/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Benoit Anastay Add-on: Paperless-ngx
# Configures Parless-ngx
# ==============================================================================

PAPERLESS_CONSUMPTION_DIR=/share/paperless/consume
PAPERLESS_DATA_DIR=/config/data
PAPERLESS_MEDIA_ROOT=/share/paperless/media

if bashio::fs.directory_exists '/share/paperless/data'; then
    bashio::log.info "Migrating 'data' folder from share storage to /addon_configs/papeless_ngx/data"
    mv /share/paperless/data /config
fi

mkdir -p $PAPERLESS_CONSUMPTION_DIR
mkdir -p $PAPERLESS_DATA_DIR
mkdir -p $PAPERLESS_MEDIA_ROOT

chown -R paperless:paperless /usr/src/paperless

cp /usr/src/paperless/paperless.conf /config/paperless.conf.exemple